# infra/ec2/systemd/clawndom-sync-agents.timer
[Unit]
Description=Periodic git pull of Clawndom agent repos

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
