# infra/ec2/systemd/clawndom.service
[Unit]
Description=Clawndom webhook proxy
After=network-online.target redis-server.service tailscaled.service
Wants=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=clawndom
Group=clawndom
WorkingDirectory=/opt/clawndom
EnvironmentFile=/etc/clawndom/clawndom.env
ExecStart=/usr/bin/node /opt/clawndom/dist/server.js
Restart=always
RestartSec=5
StandardOutput=append:/var/log/clawndom/clawndom.log
StandardError=append:/var/log/clawndom/clawndom.log

# Memory ceiling — cap the cgroup so a runaway warm-session subprocess
# (claude-cli + its MCP children) gets OOM-killed inside the unit instead
# of wedging the whole EC2 host. Sized for a t3.medium (4 GB total),
# leaving ~1 GB for the OS + redis + tailscaled + ssm-agent baseline.
# OOMPolicy=continue keeps the unit running after the kernel kills an
# offending child; Restart=always brings the unit back if the main
# process itself is the victim.
MemoryMax=3G
OOMPolicy=continue

# Hardening — narrow the blast radius if the process is compromised
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/clawndom /var/log/clawndom
PrivateTmp=true

[Install]
WantedBy=multi-user.target
