# infra/launchd/com.openclaw.jira-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.jira-proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>INSTALL_PATH/dist/server.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>JIRA_HMAC_SECRET</key>
    <string>JIRA_HMAC_SECRET</string>
    <key>OPENCLAW_TOKEN</key>
    <string>OPENCLAW_TOKEN</string>
    <key>OPENCLAW_HOOK_URL</key>
    <string>OPENCLAW_HOOK_URL</string>
    <key>REDIS_URL</key>
    <string>REDIS_URL</string>
    <key>PORT</key>
    <string>PORT</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/usr/local/var/log/openclaw-jira-proxy.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/openclaw-jira-proxy.log</string>
</dict>
</plist>
