# infra/launchd/com.openclaw.clawndom-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.clawndom-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>
PLIST="$HOME/Library/LaunchAgents/com.openclaw.clawndom.plist"
LOG="/usr/local/var/log/clawndom-refresh.log"
log() { echo "$(date '+%Y-%m-%d %H:%M:%S') $1" >> "$LOG"; }

TOKEN=$(security find-generic-password -s "Claude Code-credentials" -w 2>/dev/null | python3 -c "import sys,json; print(json.load(sys.stdin)['claudeAiOauth']['accessToken'])" 2>/dev/null)
if [ -z "$TOKEN" ]; then
  log "ERROR: Could not extract token"
  exit 1
fi

python3 -c "
import plistlib
with open('$PLIST', 'rb') as f:
    plist = plistlib.load(f)
plist['EnvironmentVariables']['CLAUDE_CODE_OAUTH_TOKEN'] = '$TOKEN'
with open('$PLIST', 'wb') as f:
    plistlib.dump(plist, f)
"
plutil -convert xml1 "$PLIST"
launchctl unload "$PLIST" 2>/dev/null
launchctl load "$PLIST"
log "OK: Token refreshed and service reloaded"
    </string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>__NODE_BIN_DIR__:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Run in the GUI session (Aqua) where Keychain is accessible -->
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>

  <!-- Run on login + every 6 hours -->
  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/usr/local/var/log/clawndom-refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/clawndom-refresh.log</string>
</dict>
</plist>
