# infra/launchd/com.openclaw.clawndom.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.clawndom</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__INSTALL_PATH__/dist/server.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>USER</key>
    <string>__USER__</string>
    <key>PATH</key>
    <string>__NODE_BIN_DIR__:/usr/local/bin:/usr/bin:/bin</string>
    <key>OP_SERVICE_ACCOUNT_TOKEN</key>
    <string>__OP_TOKEN__</string>
    <key>PROVIDERS_CONFIG</key>
    <string>__PROVIDERS_CONFIG__</string>
    <key>SECRETS_PROVIDERS_CONFIG</key>
    <string>__SECRETS_PROVIDERS_CONFIG__</string>
    <key>SECRETS_CONFIG</key>
    <string>__SECRETS_CONFIG__</string>
    <key>REDIS_URL</key>
    <string>__REDIS_URL__</string>
    <key>PORT</key>
    <string>__PORT__</string>
    <key>MAX_CONCURRENT_RUNS</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>SessionCreate</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/usr/local/var/log/clawndom.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/clawndom.log</string>
</dict>
</plist>
